# scripts/upstart.conf
start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

respawn
exec /opt/{{ cookiecutter.package_name }}/env/bin/python run.py
